# roles/cloudflare-ddns/files/cloudflare_ddns.timer
[Unit]
Description=Cloudflare DDNS timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
